# data/gmessage-daemon.service
[Unit]
Description=GMessage Background Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gmessage --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
